# services/ddclient/run
#!/bin/bash

echo "Starting ddclient"
exec /sbin/setuser abc /usr/sbin/ddclient -foreground
